# config/cron/crontabs/development.crontab
0 */2 * * * python /usr/src/app/manage.py delete_outdated_mediafiles
