# config/cron/crontabs/qa.crontab
30 4 * * * python /usr/src/app/manage.py delete_outdated_mediafiles
30 5 * * * python /usr/src/app/manage.py delete_user_authentication
